# dist/debian/mariner3d.service
[Unit]
Description=An utility for controlling Elegoo Mars 3D printers

[Service]
User=mariner
ExecStart=/usr/bin/mariner
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
